# supervisord/run-once-sample.ini
[program:zk-server]
command=/root/kafka_2.12-2.3.0/bin/zookeeper-server-start.sh /root/kafka_2.12-2.3.0/config/zookeeper.properties
autostart=true
autorestart=true
startsecs=3
priority=1

[program:kafka-server]
command=/root/kafka_2.12-2.3.0/bin/kafka-server-start.sh /root/kafka_2.12-2.3.0/config/server.properties
autostart=true
autorestart=true
startsecs=3
priority=2

[program:redis]
command=/usr/bin/redis-server
autostart=true
autorestart=true
startsecs=3
priority=3

[program:create-topic]
command=/root/kafka_2.12-2.3.0/bin/kafka-topics.sh --create --zookeeper localhost:2181 --replication-factor 1 --partitions 1 --topic WAF_COLLECTION_GZ01
autostart=true
autorestart=unexpected
priority=4
exitcodes=0
